# etc/runit/millennial-title/run
#!/bin/sh
eval $( perl -Mlocal::lib )

export MOJO_HOME=$(pwd)

exec millennial-title prefork -l 'http://*:6032' 2>&1
