# etc/runit/millennial-title/run
#!/bin/sh
eval $( perl -Mlocal::lib )

export MOJO_HOME=$(pwd)
export MOJO_REVERSE_PROXY=/title

exec millennial-title prefork -l 'http://*:6032' 2>&1
